# crackpi-client.service
[Unit]
Description=CrackPi Distributed Password Cracking Client
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/crackpi
Environment=PATH=/home/pi/crackpi/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/pi/crackpi
EnvironmentFile=-/etc/crackpi/client.env
ExecStart=/home/pi/crackpi/venv/bin/python3 /home/pi/crackpi/client.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crackpi-client

# Security settings
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/var/log/crackpi /tmp /home/pi/crackpi
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits - Allow more resources for cracking
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=95%

# Restart conditions
StartLimitInterval=300
StartLimitBurst=5

# Environment for cracking tools
Environment=HASHCAT_FORCE=1
Environment=OPENCL_FORCE_CPU=1

[Install]
WantedBy=multi-user.target
